# crontab
*/1  * * * * /usr/local/bin/ftnctl.sh
*/30 * * * * /usr/local/bin/ftnctl.sh poll
1 */12 * * * /usr/local/bin/ftnctl.sh housekeep
